# debian/postinst
#!/bin/sh

systemctl stop serial-getty@ttyS0.service
systemctl disable serial-getty@ttyS0.service
systemctl enable tunslip6.service
