# debian/postinst
#!/bin/sh

systemctl enable tunslip6.service
